# deploy/aws/arc-stack.service
[Unit]
Description=Arc agent stack (ArcUI + agents)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/arc
EnvironmentFile=/home/ubuntu/arc/.env
Environment=ARC_UI_PORT=8420
Environment=ARC_TEAM_ROOT=/home/ubuntu/arc/team
Environment=ARC_LOG_DIR=/home/ubuntu/arc/.arc-logs
Environment=ARC_TOKENS_FILE=/home/ubuntu/.arcagent/arc-stack.tokens
Environment=PATH=/home/ubuntu/arc/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/arc/scripts/arc-stack.sh start
ExecStop=/home/ubuntu/arc/scripts/arc-stack.sh stop
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
